# templates/validator/bobcat-dynconf.service.epp
[Unit]
Description=Dynamic Bobcat configuration service
After=volatilefs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dynconf_update

[Install]
WantedBy=multi-user.target
