# templates/validator/bobcat-dynconf.timer.epp
[Unit]
Description=Dynamic Bobcat configuration timer
After=network.target

[Timer]
OnCalendar=<%= $::bobcat::validator::dynconf_timer %>
Persistent=true

[Install]
WantedBy=timers.target
