# templates/validator/bobcat-dynconf.service.epp
[Unit]
Description=Dynamic Bobcat configuration service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dynconf_update
