# templates/validator/bobcat-dynconf.service.epp
[Unit]
Description=Dynamic Bobcat configuration service
After=volatilefs.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dynconf_update --fetchkdk

[Install]
WantedBy=multi-user.target
